# systemd/product-api.service
[Unit]
Description=Valio Product API
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/etymologer.com/database-backend
Environment="PATH=/var/www/etymologer.com/.venv/bin"
ExecStart=/var/www/etymologer.com/.venv/bin/python3 -m uvicorn api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/var/www/etymologer.com/logs/product_api.log
StandardError=append:/var/www/etymologer.com/logs/product_api.log

[Install]
WantedBy=multi-user.target
